# dist/launchd/com.hammer.nostromd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.hammer.nostromd</string>
  <key>ProgramArguments</key>
    <array><string>__PREFIX__/bin/nostromd</string></array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key><false/>
      <key>Crashed</key><true/>
    </dict>
  <key>StandardOutPath</key><string>__HOME__/.cache/nostromd/log/stdout.log</string>
  <key>StandardErrorPath</key><string>__HOME__/.cache/nostromd/log/stderr.log</string>
  <key>EnvironmentVariables</key>
    <dict>
      <key>RUST_LOG</key><string>info</string>
      <key>NOSTROMD_TCP_ADDR</key><string>0.0.0.0:47100</string>
      <key>HOME</key><string>__HOME__</string>
      <key>PATH</key><string>__HOME__/.local/bin:__HOME__/.cargo/bin:__HOME__/.claude/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
